# node_server.service
[Unit] 
Description=js server for CSYE-6225 
After=network.target 
[Service] 
Type=simple 
User=csye6225 
WorkingDirectory=/opt/webapp/ 
Environment=PORT=3000
Environment=HOST=localhost
Environment=USER_NAME=root
Environment=DB=demoDb
Environment=PASSWORD=Chocoslam
Environment=dialect=mysql 
ExecStart=/usr/bin/node /opt/webapp/app.js 
Restart=on-failure 
RestartSec=3
[Install] 
WantedBy=multi-user.target
